# homi/homi.service.in
[Unit]
Description=Host Orchestrated Multipath I/O

[Service]
Type=simple
StandardOutput=journal
StandardError=journal
ExecStart=@bindir@/homid

[Install]
WantedBy=multi-user.target
